# start_openvas.service
[Unit]
Description = Reverse ssh tunnel
After =  network-online.target
#Requires

[Service]
User = root
ExecStart = docker run --detach --publish 8080:9392 -e PASSWORD="0Y2KCRU4lq&u" -v /home/siteone/openvas-data:/data --name openvas immauss/openvas
RestartSec=30
Restart=always


[Install]
WantedBy = multi-user.target
