# start_vnc_session.service
[Unit]
Description = VNC Server on tcp/5901
After =  network-online.target
#Requires

[Service]
User = root
#Type = forking
ExecStart = vncserver :1
#ExecStop= vncserver -kill :1
#RestartSec=5
#Restart=always


[Install]
WantedBy = multi-user.target
