# vnc_ssh_tunnel.service
[Unit]
Description = Reverse ssh tunnel
After =  network-online.target
#Requires

[Service]
User = siteone
#Type = forking
ExecStart = ssh -vvvv -f -N -R 45555:localhost:5901 -p 45566 svc-scanner-sec@192.69.100.14
#ExecStop= /usr/bin/killall autossh
RestartSec=5
Restart=always


[Install]
WantedBy = multi-user.target
