# vnc_ssh_tunnel.service
[Unit]
Description = Reverse ssh tunnel
After =  network-online.target
#Requires

[Service]
User = siteone
#Type = forking
ExecStart = /home/$USER/Field-Vulnerability-Scanner-v2/start_tunnel.sh
ExecStop= kill -9 $(cat /home/$USER/ssh_vnc_tunnel.pid)
RestartSec=10
Restart=always


[Install]
WantedBy = multi-user.target
